# wend-solver.service
[Unit]
Description=Wend Solver — LinkedIn puzzle solver web app
After=network.target

[Service]
Type=simple
User=CHANGE_ME
WorkingDirectory=/path/to/wend-solver
ExecStart=/path/to/wend-solver/.venv/bin/python app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
